# vite-app.service
[Unit]
Description=Vite React Application
After=network.target

[Service]
Type=simple
User=homepage
Group=homepage
WorkingDirectory=/home/homepage/hello-R1
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin:/home/homepage/.local/share/pnpm
ExecStart=/usr/local/bin/pnpm run dev --host 0.0.0.0
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vite-app

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/homepage/hello-R1

[Install]
WantedBy=multi-user.target